# etc/supervisor/conf.d/20-Xvnc.conf
[program:Xvnc]
command=/usr/bin/Xvnc :0 -SecurityTypes none -AlwaysShared
autorestart = true
priority=20
