# etc/supervisor/conf.d/50-nxserver.conf
[program:nxserver]
command = /home/shakugan/nxserver.sh
startsecs = 0
autorestart = false
startretries = 1
