# etc/supervisor/conf.d/50-nomachine.conf
[program:nomachine]
priority=50
command=/bin/bash /etc/NX/nxserver --startup
command=/bin/bash tail -f /usr/NX/var/log/nxserver.log
autorestart = true
stopasgroup=true
