# etc/supervisor/conf.d/50-nxserver.conf
[program:nxserver]
command = /home/shakugan/nxserver.sh
startsecs = 0
autorestart = true
startretries = 1
user = shakugan
priority = 50
